# experiments/Set-006/X-001/templates/pytest001/etc/systemd/system/set-ipv6-netplan-once.service
[Unit]
Description=Set Netplan for IPv6 Once

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/set-ipv6-netplan-once.sh
