# experiments/Set-006/X-001/templates/pytest001/etc/systemd/system/set-ipv6-netplan-once.timer
[Unit]
Description=Look for IPv6 every 60s after guaranteed expiry (TTL)

[Timer]
OnBootSec={{ ip6_ttl }}
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
